# files/parental-privacy
#!/bin/sh /etc/rc.common

# START=95 ensures this runs after firewall and network are fully up.
# STOP=10 ensures clean teardown before network goes down on shutdown.
START=95
STOP=10

start() {
    # Re-apply SafeSearch (writes dnsmasq config and reloads dnsmasq)
    /usr/share/parental-privacy/safesearch.sh enable

    # Re-apply DoH blocking rules (nftables or iptables)
    /usr/share/parental-privacy/block-doh.sh enable

    # Re-apply bandwidth shaping if a limit is configured
    local LIMIT
    LIMIT=$(uci -q get parental_privacy.default.bandwidth)
    [ -n "$LIMIT" ] && /usr/share/parental-privacy/bandwidth.sh "$LIMIT"
}

stop() {
    /usr/share/parental-privacy/block-doh.sh disable
    /usr/share/parental-privacy/bandwidth.sh off
    /usr/share/parental-privacy/safesearch.sh disable
}

reload() {
    # Called by procd_add_reload_trigger when firewall/network/parental_privacy
    # UCI config changes — re-apply all rules so they survive a firewall restart.
    start
}
