# files/parental-privacy
#!/bin/sh /etc/rc.common

# START=95: run after firewall (START=20) and network (START=20) are up.
# STOP=89:  tear down before firewall and network go down on shutdown,
#           ensuring our nftables/tc rules are removed cleanly while the
#           interfaces they reference still exist.
START=95
STOP=89

# Expose 'reload' as a supported extra command so that
# `/etc/init.d/parental-privacy reload` works via the standard init.d
# interface, and so LuCI can call it after a firewall restart.
EXTRA_COMMANDS="reload"
EXTRA_HELP="        reload  Re-apply all parental-privacy rules (DoH, bandwidth, SafeSearch)"

start() {
    # Re-apply SafeSearch CNAME records (writes /etc/dnsmasq.d/safesearch.conf
    # and reloads dnsmasq). Safe to call when already enabled — idempotent.
    /usr/share/parental-privacy/safesearch.sh enable

    # Re-apply DoH blocking rules in nftables (fw4) or iptables (legacy).
    # block-doh.sh calls disable before enable, so this is also idempotent.
    /usr/share/parental-privacy/block-doh.sh enable

    # Re-apply bandwidth shaping via tc if a limit is stored in UCI.
    # If no limit is set we leave tc untouched (unlimited).
    local LIMIT
    LIMIT=$(uci -q get parental_privacy.default.bandwidth)
    [ -n "$LIMIT" ] && /usr/share/parental-privacy/bandwidth.sh "$LIMIT"
}

stop() {
    # Remove DoH blocking rules before the firewall/network interfaces go away.
    /usr/share/parental-privacy/block-doh.sh disable

    # Remove tc bandwidth shaping qdiscs.
    /usr/share/parental-privacy/bandwidth.sh off

    # Remove SafeSearch dnsmasq config and reload dnsmasq.
    /usr/share/parental-privacy/safesearch.sh disable
}

reload() {
    # Re-apply all rules without a full stop/start cycle.
    # Called after a firewall restart (which flushes our nftables rules)
    # or when UCI config changes.
    start
}
